# templates/demoapp.conf.j2
[program:demoapp]
command = {{ venv_path }}/bin/gunicorn_start
directory = {{ project_path }}
user = vagrant
autostart = true
autorestart = true
redirect_stderr = True
